# home/Library/LaunchAgents/dev.remote.postgres.mise.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" 
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>dev.remote.postgres.mise</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/brunno.dossantosleite/.local/share/mise/installs/postgres/17.5/bin/postgres</string>
    <string>-D</string>
    <string>/Users/brunno.dossantosleite/.local/share/mise/installs/postgres/17.5/data</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/tmp/postgres.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/postgres.err</string>
</dict>
</plist>
